# site-cookbooks/tmate/templates/default/tmate.conf.erb
description "tmate"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10000 30
kill timeout 5

setuid root
setgid root

script
  cd <%= @app_path %>
  echo $$ > /var/run/tmate.pid
  exec ./tmate-slave -l <%= @log_path %>/tmate.log -k <%= @keys_dir %> -h <%= @hostname %> -p <%= @listen_port %> -v
end script

post-stop script
  rm -f /var/run/tmate.pid
end script
